# board/microchip/sama5/fs-overlay/usr/lib/systemd/system/hostapd@.service.example
[Unit]
Description=Advanced IEEE 802.11 AP and IEEE 802.1X/WPA/WPA2/EAP Authenticator
Before=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/hostapd /etc/wilc_hostapd_%I.conf
